# .dist/linux/usr/share/eg/rootfs/etc/profile.d/eg-completions.sh
#!/bin/sh

source /usr/share/bash-completion/completions/eg-daemon.sh