# .dist/linux/etc/profile.d/eg.completions.sh
#!/bin/sh

source /usr/share/bash-completion/completions/eg.sh
