# .dist/linux/usr/lib/systemd/system/eg-schedule-15m.target
[Unit]
Description=for running tasks once every 15 minutes.
After=default.target
StopWhenUnneeded=true

[Install]
WantedBy=default.target
